# playbooks/roles/graphite/templates/graphite-api/systemd/graphite-api.socket.j2
[Unit]
Description=graphite-api socket

[Socket]
SocketUser=graphite
SocketGroup=graphite
ListenStream=/opt/graphite/run/graphite-api.sock

[Install]
WantedBy=sockets.target
